# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:netease-api]
command=node /app/netease-api/app.js
directory=/app/netease-api
autostart=true
autorestart=true
stdout_logfile=/var/log/netease-api.log
stderr_logfile=/var/log/netease-api-error.log

[program:web-server]
command=/app/endfield-music-server --config /data/config/config.yaml
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/web-server.log
stderr_logfile=/var/log/web-server-error.log
